# deploy/systemd/qsz-academic-deploy.timer
[Unit]
Description=Run qsz_academic deployment timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Persistent=true
Unit=qsz-academic-deploy.service

[Install]
WantedBy=timers.target
